# buildroot-external/rootfs-overlay/usr/lib/systemd/system/oppos-persists.service
[Unit]
Description=OppOS persists data on bootstate 
After=sysinit.target
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/libexec/oppos-persists
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
