# buildroot-external/rootfs-overlay/usr/lib/systemd/system/zram-var.service
[Unit]
Description=OppOS ZRAM var
DefaultDependencies=no
Before=var.mount
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/libexec/oppos-zram -t fs -m var -s 32M
RemainAfterExit=true
StandardOutput=null

[Install]
WantedBy=oppos-zram.target
