# buildroot-external/rootfs-overlay/etc/systemd/system/getty@.service.d/oppos.conf
[Unit]
Wants=oppos-supervisor.service
After=oppos-supervisor.service
